# laradock/workspace/crontab/laradock
*/1 * * * * node /var/www/scripts/sync.js index update > /dev/null 2>&1
*/2 * * * * node /var/www/scripts/sync.js market > /dev/null 2>&1
*/5 * * * * node /var/www/scripts/peers.js > /dev/null 2>&1